# webapp/mywebbapp.service
[Unit]
Description=uWSGI instance to serve mywebbapp
After=network.target
